# freeswitch.service
[Service]
Type=forking
ExecStart=/usr/bin/freeswitch -ncwait -nonat
